# packages/sysutils/busybox/system.d/etc.mount
[Unit]
Description=etc directory overlay
After=storage-dirs.service
Requires=storage-dirs.service
ConditionPathIsDirectory=/storage/.filesystem/etc
ConditionPathIsDirectory=/storage/.overlays/etc

[Mount]
What=overlay
Where=/etc
Type=overlay
Options=lowerdir=/etc,upperdir=/storage/.filesystem/etc,workdir=/storage/.overlays/etc

[Install]
WantedBy=local-fs.target
