# packages/sysutils/busybox/system.d/storage-dirs.service
[Unit]
Description=Create directories in /storage

[Service]
Type=oneshot
ExecStartPre=mkdir -p /storage/.filesystem/etc
ExecStartPre=mkdir -p /storage/.overlays/etc
ExecStart=/usr/bin/true
RemainAfterExit=yes

[Install]
WantedBy=basic.target
