# packages/sysutils/busybox/system.d/var-lib-extensions.mount
[Unit]
Description=systemd-sysext enablement
RequiresMountsFor=/var /storage
ConditionPathExists=|/storage/.config/extensions

[Mount]
What=/storage/.config/extensions
Where=/var/lib/extensions
Options=bind

[Install]
WantedBy=local-fs.target
