# packages/sysutils/entropy/system.d/add-entropy.service
[Unit]
Description=Add random entropy from file
DefaultDependencies=no
RequiresMountsFor=/storage/cache/random.data
After=local-fs.target
Requires=local-fs.target
Before=sysinit.target
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/lib/entropy/add-entropy
ExecStop=/usr/lib/entropy/add-random-at-shutdown
RemainAfterExit=yes

[Install]
WantedBy=basic.target
